# com.user.cursor-journal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.cursor-journal</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/rorymcdaniel/.local/bin/cursor-journal/daily-journal.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Backup: also run every 6 hours as catch-up for missed schedules -->
    <key>StartInterval</key>
    <integer>21600</integer>
    
    <key>StandardOutPath</key>
    <string>/tmp/cursor-journal.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/cursor-journal.stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/rorymcdaniel/.local/bin</string>
        <key>HOME</key>
        <string>/Users/rorymcdaniel</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
